# docker/root/etc/s6-overlay/s6-rc.d/init-chown/run
#!/usr/bin/with-contenv bash

echo "-------------------------------------"
# permissions
echo "chown'ing directory to ensure correct permissions."
chown -R abc:abc /home/abc
chown -R abc:abc /config
echo "Done!"
echo -e "-------------------------------------\n"
