# docker/root/etc/s6-overlay/s6-rc.d/init-chown/run
#!/usr/bin/with-contenv bash

echo "-------------------------------------"
# permissions
#echo "chown'ing home directory to ensure correct permissions."
chown -R abc:abc /home/abc
#echo "Done!"
echo -e "-------------------------------------\n"
